# platforms/cto-ai/supervisor/browser-health.conf
[program:browser-health]
command=/bin/bash -c "sleep 10 && while true; do /root/health-check.sh; sleep 300; done"
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/browser-health.log
stderr_logfile=/var/log/browser-health.err
environment=TARGET_URL="%(ENV_TARGET_URL)s",VNC_PORT="%(ENV_VNC_PORT)s",ALERT_WEBHOOK="%(ENV_ALERT_WEBHOOK)s"
