# platforms/cto-ai/supervisor/browser-novnc.conf
[program:browser-novnc]
command=novnc_proxy --vnc localhost:5900 --listen 3000
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/browser-novnc.log
stderr_logfile=/var/log/browser-novnc.err
environment=DISPLAY=":99"
